# templates/init-script.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides:          bamboo-agent<%= @agent_id.to_s.ljust(4) %>
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Init script for running bamboo agent
### END INIT INFO
# bamboo-agent<%= @agent_id.to_s.ljust(4) %> Init script for running bamboo agent
#
# chkconfig:       2345 98 02
#
# description:     Starts and stops bamboo agent <%= @agent_id %>

# This is just a delegator to the Bamboo agent script.

if [[ "$EUID" != '0' ]]; then
  echo 'This script must be run as root!'
  exit 1
fi

status=`su - <%= @user %> -c "<%= @script %> $@ 1>&2 ;"' echo $?'`
exit $status
